# tools/systemd/njcic-explorer-refresh.service
[Unit]
Description=Refresh NJCIC grantee explorer data from Airtable
After=network-online.target
Wants=network-online.target
# Don't swamp the journal if something goes wrong and it fails repeatedly.
# (These directives only take effect in [Unit] — systemd silently ignores
# them when placed under [Service].)
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=oneshot
User=jamditis
Group=jamditis
# Refuse to run if a previous invocation (timer-driven OR hand-run with
# any args) is still going. flock(1) catches all callers because they
# all contend for the same lockfile, unlike pgrep -fx which only matches
# the exact argv we expect from the timer.
ExecStart=/usr/bin/flock -n /tmp/njcic-explorer-refresh.lock /bin/bash /home/jamditis/projects/njcic/tools/refresh_explorer.sh
# Allow sudo cp (already passwordless on this host per CLAUDE.md)
StandardOutput=journal
StandardError=journal
